# samknows_base-files-router/files/etc/init.d/pm
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/lib/samknows/setup.sh > /dev/null 2>&1 &
}
