# samknows_base-files-router/files/etc/rc.local
crontab /etc/crontabs/mycrontab 2>/dev/null
bwmonitor &
#start editing now

exit 0